# lab9/mephi.service
[Unit]
Decription=test
After=sshd.service
[Service]
ExecStart=/usr/bin/printf "hello world"
[Install]
WantedBy=multi-user.target
